# installer/ubuntu/package/DEBIAN/postinst
#!/bin/sh

# エラーが発生した時点で処理を終了する
set -e

###############################################################################
# MonoMaint Debian Package Post Install Script
#
# Debian パッケージインストール完了後に実行されます。
#
# 実施内容
#   1. サービス実行用グループ作成
#   2. サービス実行用ユーザー作成
#   3. ログディレクトリ作成
#   4. アプリケーション所有者設定
#   5. systemd再読み込み
#   6. サービス有効化
#   7. サービス起動
###############################################################################

# systemdサービス実行ユーザー
SERVICE_USER="monomaint"

# systemdサービス実行グループ
SERVICE_GROUP="monomaint"

# アプリケーション配置先
APP_DIR="/opt/monomaint"

# ログ出力先
LOG_DIR="/var/log/monomaint"

# systemdサービス名
SERVICE_NAME="monomaint-sample.service"


###############################################################################
# サービス実行グループ作成
###############################################################################

echo "Creating service group: ${SERVICE_GROUP}"

if ! getent group "${SERVICE_GROUP}" >/dev/null 2>&1; then
    groupadd --system "${SERVICE_GROUP}"
fi


###############################################################################
# サービス実行ユーザー作成
###############################################################################

echo "Creating service user: ${SERVICE_USER}"

if ! getent passwd "${SERVICE_USER}" >/dev/null 2>&1; then
    useradd \
        --system \
        --gid "${SERVICE_GROUP}" \
        --home-dir "${APP_DIR}" \
        --shell /usr/sbin/nologin \
        "${SERVICE_USER}"
fi


###############################################################################
# ログディレクトリ作成
#
# Serilogなどのログ出力先
###############################################################################

echo "Creating log directory..."

mkdir -p "${LOG_DIR}"

chown "${SERVICE_USER}:${SERVICE_GROUP}" "${LOG_DIR}"

chmod 0750 "${LOG_DIR}"


###############################################################################
# アプリケーション配置先の所有者変更
###############################################################################

echo "Setting application ownership..."

chown -R "${SERVICE_USER}:${SERVICE_GROUP}" "${APP_DIR}"


###############################################################################
# systemd設定反映
#
# Unitファイルを再読み込みし、
# OS起動時に自動起動するよう登録したあと、
# サービスを起動します。
###############################################################################

echo "Reloading systemd..."

systemctl daemon-reload

systemctl enable "${SERVICE_NAME}"

systemctl restart "${SERVICE_NAME}"


###############################################################################
# 完了
###############################################################################

echo "MonoMaint installation completed."

exit 0